# varnish/files/supervisor.conf
{%- set lookup = salt['pillar.get']('varnish:server:lookup:'+service_name) %}

[program:varnish_{{ lookup.name }}]
command=varnishd -F -P /var/run/varnish.pid -a {{ lookup.bind.host }}:{{ lookup.bind.port }} -f /etc/varnish/{{ lookup.name }}.vcl -T localhost:6082 -t 120 -u varnish
stdout_logfile=/var/log/varnish/{{ lookup.name }}_access.log 
stderr_logfile=/var/log/varnish/{{ lookup.name }}_error.log
user=root
autostart=true
autorestart=true
redirect_stderr = true